# roles/setup_nfs_ldap/templates/users-home.automount
[Unit]
Description=Automount Home

[Automount]
Where={{ nfs_home_path }}

[Install]
WantedBy=multi-user.target
